# ops/com.auto-browser.marketplace-db-maintenance.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.auto-browser.marketplace-db-maintenance</string>

  <key>WorkingDirectory</key>
  <string>/path/to/auto-browser</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>npm</string>
    <string>run</string>
    <string>marketplace:home:db:maintain</string>
    <string>--</string>
    <string>--json</string>
  </array>

  <key>StartInterval</key>
  <integer>21600</integer>

  <key>StandardOutPath</key>
  <string>/path/to/auto-browser/artifacts/marketplace-homepage/db-maintenance.log</string>

  <key>StandardErrorPath</key>
  <string>/path/to/auto-browser/artifacts/marketplace-homepage/db-maintenance.err.log</string>
</dict>
</plist>
