# ops/marketplace-db-analyze.timer.example
[Unit]
Description=Run Marketplace homepage SQLite statistics refresh weekly

[Timer]
OnCalendar=Sun 03:30
Persistent=true
Unit=marketplace-db-analyze.service

[Install]
WantedBy=timers.target
