# ops/marketplace-db-analyze.service.example
[Unit]
Description=Marketplace homepage SQLite weekly statistics refresh
Documentation=README.md

[Service]
Type=oneshot
User=auto-browser
Group=auto-browser
WorkingDirectory=/srv/auto-browser/app
ExecStart=/usr/bin/env npm run marketplace:home:db:maintain -- --analyze --json
StandardOutput=append:/srv/auto-browser/app/artifacts/marketplace-homepage/db-maintenance.log
StandardError=append:/srv/auto-browser/app/artifacts/marketplace-homepage/db-maintenance.err.log
